# image/etc/rc.local
#!/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
# set hostname
hostname ‘/bin/whereami‘
exit 0
